# supervisord.conf
[supervisord]
nodaemon=true

[program:X11]
command=Xvfb :0 -screen 0 "%(ENV_DISPLAY_WIDTH)s"x"%(ENV_DISPLAY_HEIGHT)s"x24
autorestart=true

[program:x11vnc]
command=x11vnc -forever
autorestart=true

[program:novnc]
command=websockify --web /usr/share/novnc 8080 localhost:5900
autorestart=true

[program:socat]
command=socat tcp-listen:6000,reuseaddr,fork unix:/tmp/.X11-unix/X0
autorestart=true

[program:fluxbox]
command=fluxbox
autorestart=true

[program:demo]
command=xterm
autorestart=true
