# conf.d/demo.conf
[program:demo]
command=xterm
autorestart=true
