# conf.d/websockify.conf
[program:websockify]
command=websockify --web /usr/share/novnc 8081 localhost:5900
autorestart=true
